# pkg/deploy/config/ograc_logs_handler.timer
[Unit]
Description=Run every 60minutes and on boot

[Timer]
OnBootSec=2min
OnUnitActiveSec=60min

[Install]
WantedBy=multi-user.target